# deploy/ksm/redtusk-ksm.service
[Unit]
Description=Enable KSM for RedTusk worker memory dedup
After=network.target
ConditionPathExists=/sys/kernel/mm/ksm/run

[Service]
Type=oneshot
ExecStart=/usr/local/bin/redtusk-enable-ksm.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
